# config/systemd/jericho-host-bridge.service.template
[Unit]
Description=Jericho Host Bridge (Kimi Web UI Spawner)
After=network.target

[Service]
Type=simple
WorkingDirectory=${INSTALL_DIR}
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 ${INSTALL_DIR}/host-bridge.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
